# systemd/zeicoin-rest-api.service
[Unit]
Description=ZeiCoin Transaction API Service (Hybrid RPC + PostgreSQL) - Legacy Alias
After=zeicoin-server.service postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/zeicoin
EnvironmentFile=-/root/zeicoin/.env
Environment="ZEICOIN_SERVER=127.0.0.1"
ExecStart=/root/zeicoin/zig-out/bin/transaction_api
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target