# systemd/zeicoin-transaction-api.service
# SPDX-FileCopyrightText: 2025-2026 Ryo Zen (https://github.com/ryo-zen)
# SPDX-License-Identifier: Apache-2.0

[Unit]
Description=ZeiCoin Transaction API Service (Hybrid RPC + PostgreSQL)
After=zeicoin-mining.service postgresql.service
Wants=zeicoin-mining.service postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/zeicoin
EnvironmentFile=-/root/zeicoin/.env.testnet
EnvironmentFile=-/root/zeicoin/.env.local
Environment="ZEICOIN_SERVER=127.0.0.1"

# Clean up stale api processes
ExecStartPre=-/usr/bin/pkill -9 -x transaction_api

ExecStart=/root/zeicoin/zig-out/bin/transaction_api
Restart=always
RestartSec=5
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
